# .config/autostart/clipmenud.desktop
[Desktop Entry]
Name=Clipmenud
Type=Application
Exec=/usr/bin/clipmenud
Terminal=false
