# .config/autostart/keyboard-center.desktop
[Desktop Entry]
Name=Keyboard Center
Type=Application
Exec=/usr/bin/keyboard-center
Icon=keyboard-center
Hidden=true
